# yaskkserv-git/yaskkserv.socket
[Unit]
Description=yaskkserv socket
PartOf=yaskkserv.service

[Socket]
ListenStream=11780

[Install]
WantedBy=sockets.target
